# deb/DEBIAN/postrm
#!/bin/bash
set -e

# Restore original proxychains4.conf if backup exists
if [ -f /etc/loupvagabond/proxychains4.conf.backup ]; then
    echo "Restoring original proxychains4.conf from backup..."
    cp /etc/loupvagabond/proxychains4.conf.backup /etc/proxychains4.conf
    
    # Remove backup if requested to purge
    if [ "$1" = "purge" ]; then
        rm -f /etc/loupvagabond/proxychains4.conf.backup
    fi
fi

# Note about Tor service
if systemctl is-enabled --quiet tor.service 2>/dev/null; then
    echo "Note: The Tor service remains enabled. If you no longer need it,"
    echo "you can disable it with: sudo systemctl disable tor.service"
    echo "or remove it completely with: sudo apt remove tor"
fi

# Update desktop database
if [ -x "$(command -v update-desktop-database)" ]; then
    update-desktop-database -q
fi

# Remove config files if purge was requested
if [ "$1" = "purge" ]; then
    # Remove any configuration files or leftover data
    echo "Purging all LoupVagabond data..."
    rm -rf /etc/loupvagabond
fi

exit 0
